# Docker/preload/crontab
0 * * * * /app/preload_cache.sh >> /app/preload_log.txt 2>&1
